# scheduler/crontab
0 * * * * cd /app && python -m pipeline.runner >> /var/log/cron.log 2>&1
